# lora2mqtt/rootfs/etc/services.d/lora2mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the lora2mqtt service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
#declare loglevel
declare mqttbrokerhost
declare mqttbrokerport
declare mqttbrokeruser
declare mqttbrokerpass

# Get the keys from the user config options.
#log_level=$(bashio::config 'log_level')

# Get the keys from the mqtt service.
mqttbrokerhost=$(bashio::services mqtt "host")
mqttbrokerport=1883
mqttbrokeruser=$(bashio::services mqtt "username")
mqttbrokerpass=$(bashio::services mqtt "password")

# Print
bashio::log.info "Run application..."
bashio::log.info mqttbrokerhost
bashio::log.info mqttbrokerport
bashio::log.info mqttbrokeruser
bashio::log.info mqttbrokerpass

# Run application
exec python3 /usr/bin/bridge.py --broker=${mqttbrokerhost} --port=${mqttbrokerport} --user=${mqttbrokeruser} --Pass=${mqttbrokerpass}
