# lora2mqtt/rootfs/etc/services.d/lora2mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the lora2mqtt service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare home_assistant_prefix
declare chip_mac
declare lora_slave_addrs
declare lora_slave_names
declare lora_slave_macs
declare lora_slave_vers
declare lora_slave_chips
declare loglevel
declare maxthreads
declare mqttbrokerhost
declare mqttbrokerport
declare mqttbrokeruser
declare mqttbrokerpass

# Get the keys from the user config options.
home_assistant_prefix=$(bashio::config 'ttlockclientid')
chip_mac=$(bashio::config 'ttlockclientsecret')
lora_slave_addrs=$(bashio::config 'ttlockuser')
lora_slave_names=$(bashio::config 'ttlockhash')
lora_slave_macs=$(bashio::config 'publishstatedelay')
lora_slave_vers=$(bashio::config 'publishstatedelayshort')
lora_slave_chips=$(bashio::config 'publishbatterydelay')
loglevel=$(bashio::config 'loglevel')
maxthreads=$(bashio::config 'maxthreads')

# Get the keys from the mqtt service.
mqttbrokerhost=$(bashio::services mqtt "host")
mqttbrokerport=1883
mqttbrokeruser=$(bashio::services mqtt "username")
mqttbrokerpass=$(bashio::services mqtt "password")

# Print
bashio::log.info "Run application..."

# Run application
exec python3 /usr/bin/lora_adapter.py --haprefix=${home_assistant_prefix} --chip=${chip_mac} --addrs=${lora_slave_addrs} --names=${lora_slave_names} --broker=${mqttbrokerhost} --port=${mqttbrokerport} --user=${mqttbrokeruser} --Pass=${mqttbrokerpass} --macs=${lora_slave_macs} --vers=${lora_slave_vers} --Chips=${lora_slave_chips} --log_level=${loglevel} --Max_threads=${maxthreads}
