# lora2mqtt/rootfs/etc/services.d/lora2mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the lora2mqtt service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare homeassistantprefix
declare chipmac
declare loraslaveaddrs
declare loraslavenames
declare loraslavemacs
declare loraslavevers
declare loraslavechips
declare loglevel
declare maxthreads
declare mqttbrokerhost
declare mqttbrokerport
declare mqttbrokeruser
declare mqttbrokerpass

# Get the keys from the user config options.
homeassistantprefix=$(bashio::config 'homeassistantprefix')
chipmac=$(bashio::config 'chipmac')
loraslaveaddrs=$(bashio::config 'loraslaveaddrs')
loraslavenames=$(bashio::config 'loraslavenames')
loraslavemacs=$(bashio::config 'loraslavemacs')
loraslavevers=$(bashio::config 'loraslavevers')
loraslavechips=$(bashio::config 'loraslavechips')
loglevel=$(bashio::config 'loglevel')
maxthreads=$(bashio::config 'maxthreads')

# Get the keys from the mqtt service.
mqttbrokerhost=$(bashio::services mqtt "host")
mqttbrokerport=1883
mqttbrokeruser=$(bashio::services mqtt "username")
mqttbrokerpass=$(bashio::services mqtt "password")

# Print
bashio::log.info "Run application..."
bashio::log.info $(bashio::services mqtt "host")
bashio::log.info 1883
bashio::log.info $(bashio::services mqtt "username")
bashio::log.info $(bashio::services mqtt "password")

# Run application
exec python3 /usr/bin/lora_adapter.py --haprefix=${homeassistantprefix} --chip=${chipmac} --addrs=${loraslaveaddrs} --names=${loraslavenames} --broker=${mqttbrokerhost} --port=${mqttbrokerport} --user=${mqttbrokeruser} --Pass=${mqttbrokerpass} --macs=${loraslavemacs} --vers=${loraslavevers} --Chips=${loraslavechips} --log_level=${loglevel} --Max_threads=${maxthreads}
